# vercel.cron
# vercel.cron
0 */2 * * * api/cron/finalize-deletes
